# bluetooth/templates/bluetooth-connection.service.template
[Unit]
Description=Bluetooth Connection and Data Monitoring
After=bluetooth.service

[Service]
ExecStart=$MONITOR

[Install]
WantedBy=multi-user.target
